# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

# FastAPI Web Server
[program:api]
command=python -m uvicorn api.app:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/api.err.log
stdout_logfile=/var/log/api.out.log
environment=PATH="/app:%(ENV_PATH)s",DATABASE_URL="%(ENV_DATABASE_URL)s",TELEGRAM_BOT_TOKEN="%(ENV_TELEGRAM_BOT_TOKEN)s",GOOGLE_DRIVE_ENABLED="%(ENV_GOOGLE_DRIVE_ENABLED)s",GOOGLE_DRIVE_CREDENTIALS_PATH="%(ENV_GOOGLE_DRIVE_CREDENTIALS_PATH)s",GOOGLE_DRIVE_FOLDER_ID="%(ENV_GOOGLE_DRIVE_FOLDER_ID)s",AUDIO_DIR="/app/audio",TEMP_AUDIO_DIR="/app/temp_audio"

# Telegram Bot Worker
[program:bot]
command=python main.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/bot.err.log
stdout_logfile=/var/log/bot.out.log
environment=PATH="/app:%(ENV_PATH)s",DATABASE_URL="%(ENV_DATABASE_URL)s",TELEGRAM_BOT_TOKEN="%(ENV_TELEGRAM_BOT_TOKEN)s",GOOGLE_DRIVE_ENABLED="%(ENV_GOOGLE_DRIVE_ENABLED)s",GOOGLE_DRIVE_CREDENTIALS_PATH="%(ENV_GOOGLE_DRIVE_CREDENTIALS_PATH)s",GOOGLE_DRIVE_FOLDER_ID="%(ENV_GOOGLE_DRIVE_FOLDER_ID)s",AUDIO_DIR="/app/audio",TEMP_AUDIO_DIR="/app/temp_audio"

# Health Check (prevents Render free tier spindown)
[program:healthcheck]
command=python healthcheck.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/healthcheck.err.log
stdout_logfile=/var/log/healthcheck.out.log

# HTTP server for supervisord web interface (optional)
[inet_http_server]
port=0.0.0.0:9001
